# swmpjs.service
[Unit]
Description=SWMPjs

[Service]
Type=simple
ExecStart=/usr/local/share/swmpjs/run.sh

[Install]
WantedBy=multi-user.target
